# SPECS/acpid/acpid.path
[Unit]
Description=ACPI Events Check
ConditionVirtualization=!container

[Path]
DirectoryNotEmpty=/etc/acpi/events/

[Install]
WantedBy=paths.target

